# jobs/dellog.cron
*/12 * * * * /root/gem/scripts/dellog.sh
