# jobs/dellog.cron
*/5 * * * * /root/gem/scripts/dellog.sh
