# jobs/taskverify.cron
0 7 * * * curl http://0.0.0.0:7443/bot/check-users-tasks
